# config/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/sbin/php-fpm5.6

[program:nginx]
command=/usr/sbin/nginx